# sunrisewinds.service
[Unit]
Description=Sunrise Winds Aurora Monitor
After=network.target

[Service]
Type=simple
User=sunrisewinds
WorkingDirectory=/opt/sunrisewinds
ExecStart=/opt/sunrisewinds/sunrisewinds
Restart=always
RestartSec=5
Environment=RUST_LOG=sunrisewinds=info

[Install]
WantedBy=multi-user.target
